# sbb-flip-clock.service
[Unit]
Description=SBB Flip Clock

[Service]
ExecStart=/home/sbb/.local/bin/uv run python /home/sbb/Desktop/sbb-flip-clock/main.py
WorkingDirectory=/home/sbb/Desktop/sbb-flip-clock
Restart=always
RestartSec=5
StartLimitBurst=0

[Install]
WantedBy=default.target
